# contestImage/v1.1/casper-rw/upper/etc/rc.local
rm -f /var/lib/snapd/system-key
systemctl restart snapd
